# deploy/ai-secretary.service
[Unit]
Description=AI Secretary FastAPI service (secr-app-d01)
Documentation=file:///opt/ai-secretary/INSTALL.md
After=network-online.target docker.service ollama.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=aisec
Group=aisec
WorkingDirectory=/opt/ai-secretary
EnvironmentFile=/opt/ai-secretary/.env

# Pre-start: bring up dockerized dependencies and run migrations.
# Compose is idempotent — repeated `up -d` is a no-op when containers are healthy.
ExecStartPre=/usr/bin/docker compose -f /opt/ai-secretary/docker-compose.yml --env-file /opt/ai-secretary/.env up -d
ExecStartPre=/opt/ai-secretary/.venv/bin/alembic -c /opt/ai-secretary/alembic.ini upgrade head

ExecStart=/opt/ai-secretary/.venv/bin/uvicorn app.main:app \
  --host ${APP_HOST} \
  --port ${APP_PORT} \
  --workers ${APP_WORKERS} \
  --timeout-keep-alive 75 \
  --proxy-headers \
  --forwarded-allow-ips=127.0.0.1

Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=60

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

# Whisper / ML libs occasionally need scratch space
ReadWritePaths=/opt/ai-secretary /tmp

[Install]
WantedBy=multi-user.target
